# contrib/init/bitcoind.conf
description "FirstCoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env FIRSTCOIND_BIN="/usr/bin/firstcoind"
env FIRSTCOIND_USER="firstcoin"
env FIRSTCOIND_GROUP="firstcoin"
env FIRSTCOIND_PIDDIR="/var/run/firstcoind"
# upstart can't handle variables constructed with other variables
env FIRSTCOIND_PIDFILE="/var/run/firstcoind/firstcoind.pid"
env FIRSTCOIND_CONFIGFILE="/etc/firstcoin/firstcoin.conf"
env FIRSTCOIND_DATADIR="/var/lib/firstcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # firstcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$FIRSTCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run firstcoind."
        echo "The setting must appear in $FIRSTCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"FirstCoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$FIRSTCOIND_PIDDIR"
    chmod 0755 "$FIRSTCOIND_PIDDIR"
    chown $FIRSTCOIND_USER:$FIRSTCOIND_GROUP "$FIRSTCOIND_PIDDIR"
    chown $FIRSTCOIND_USER:$FIRSTCOIND_GROUP "$FIRSTCOIND_CONFIGFILE"
    chmod 0660 "$FIRSTCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$FIRSTCOIND_PIDFILE" \
    --chuid $FIRSTCOIND_USER:$FIRSTCOIND_GROUP \
    --exec "$FIRSTCOIND_BIN" \
    -- \
    -pid="$FIRSTCOIND_PIDFILE" \
    -conf="$FIRSTCOIND_CONFIGFILE" \
    -datadir="$FIRSTCOIND_DATADIR" \
    -disablewallet \
    -daemon

